# supervisor/camera.conf
[program:camera]
command=bash -c "/workspaces/Robot_Pi_4WD_ROS/start_camera.sh"
