# supervisor/depth.conf
[program:depth]
command=bash -c "python3 /workspaces/Robot_Pi_4WD_ROS/tof_publisher.py"
