# supervisor/imu2.conf
[program:imu2]
command=bash -c "/workspaces/Robot_Pi_4WD_ROS/start_bno055_imu.sh"
stopsignal=KILL
killasgroup=true
stopasgroup=true