# supervisor/microros.conf
[program:microros]
command=bash -c "/workspaces/Robot_Pi_4WD_ROS/start_microros.sh"