# supervisor/speech.conf
[program:speech]
command=bash -c "python3 /workspaces/Robot_Pi_4WD_ROS/speech-recognition/listen.py"