# pi-files/camera.conf
[program:camera]
command=bash -c "~/Robot_Pi_4WD_ROS/pi-files/start_camera.sh"
stopsignal=KILL
killasgroup=true
stopasgroup=true
